# gunicorn_birdnest.socket
[Unit]
Description=gunicorn birdnest socket

[Socket]
ListenStream=/run/gunicorn_birdnest.sock

[Install]
WantedBy=sockets.target

